# deploy/noteweaver.service
[Unit]
Description=NoteWeaver Knowledge Harness Gateway
After=network.target

[Service]
Type=simple
User=noteweaver
WorkingDirectory=/home/noteweaver
EnvironmentFile=/home/noteweaver/.noteweaver.env
Environment=NW_VAULT=/home/noteweaver/vault
ExecStart=/home/noteweaver/.nw-venv/bin/nw gateway
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
